# luci-app-xkcptun/root/etc/uci-defaults/90_luci_xkcptun
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@xkcptun[-1]
	add ucitrack xkcptun
	set ucitrack.@xkcptun[-1].init=kcptun_client
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0