# dnspod-utils/files/etc/hotplug.d/iface/80-dnspod
#!/bin/sh

if [ "$INTERFACE" = "wan" -a "$ACTION" = ifup ]; then
	start-stop-daemon -S -b -x /usr/bin/dnspod-utils.sh -- -d 5 -k once
fi

